# deploy/outline-webapp.service
[Unit]
Description=Outline Control Room (web dashboard)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# adjust User and paths to your setup
User=outline
WorkingDirectory=/opt/outline
EnvironmentFile=/opt/outline/.env
ExecStart=/opt/outline/.venv/bin/outline-panel
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
