# deploy/outline-panel.service
[Unit]
Description=Outline Panel (web dashboard)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# adjust User and paths to your setup
User=outline
WorkingDirectory=/opt/outline-panel
EnvironmentFile=/opt/outline-panel/.env
ExecStart=/opt/outline-panel/.venv/bin/outline-panel
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
